# deploy/systemd/janushc-recording-ingest.timer
[Unit]
Description=Run JanusHC recording ingest every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
RandomizedDelaySec=2min
Persistent=true

[Install]
WantedBy=timers.target
